# client.example.com/etc/network/if-up.d/ssh-tunnel
#!/bin/bash
ssh shell.example.com -F /etc/ssh/openvpn_tunnel_config -N &
